# deploy/supervisor/cora-server.conf
[program:cora-server]
command=/home/cora/cora-server -addr=10.0.0.10:8080 -db=/home/cora/cora.db -auth-token-file=/home/cora/auth.token
directory=/home/cora
user=cora
autostart=true
autorestart=unexpected
exitcodes=0
startsecs=3
startretries=5
stopsignal=TERM
stopwaitsecs=15
stopasgroup=true
killasgroup=true
redirect_stderr=true
stdout_logfile=/home/cora/cora-server.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=5
